# paperwallet/paperwallet/rpi/udev/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/paperwallet/rpi/bin/starter.sh"
